# launchd/com.mayflies.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template only — setup.py generates and installs the real version
  with correct absolute paths for your machine.

  Runs once daily at 09:00. Adjust Hour/Minute to your preference.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mayflies.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/uv</string>
        <string>run</string>
        <string>/path/to/mayflies/daily.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/path/to/mayflies/logs/daily.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/mayflies/logs/daily.err</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
