# launchd/com.mayflies.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template only — setup.py generates and installs the real version
  with correct absolute paths for your machine.

  To install manually, replace /path/to/mayflies with your actual path,
  copy to ~/Library/LaunchAgents/, then:
    launchctl load ~/Library/LaunchAgents/com.mayflies.poll.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mayflies.poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/uv</string>
        <string>run</string>
        <string>/path/to/mayflies/poll.py</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/path/to/mayflies/logs/poll.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/mayflies/logs/poll.err</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
